# src/docker/buildfiles/etc/supervisor.d/003-app-bootstrap.ini
[program:app-create-db]
command=sh -c 'cd /opt/ada/app && php bin/console doctrine:database:create --no-interaction --no-ansi'
priority=2
user=ada
autorestart=unexpected
exitcodes=0
startretries=3
startsecs=0
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:app-migrate-db]
command=sh -c 'cd /opt/ada/app && php bin/console doctrine:migrations:migrate --no-interaction --no-ansi'
priority=3
user=ada
autorestart=unexpected
exitcodes=0
startretries=3
startsecs=0
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:app-start-watcher]
command=sh -c 'cd /opt/ada/app && php bin/console app:watcher'
priority=4
autorestart=true
user=ada
startretries=2147483647
startsecs=5
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
